# waptserver/scripts/waptserver.service
[Unit]
Description=WAPT Server startup script
After=syslog.target
After=postgresql-9.4.service


[Service]
Type=simple
User=wapt
Group=www-data
WorkingDirectory=/opt/wapt
ExecStart=/usr/bin/python /opt/wapt/waptserver/waptserver.py
Restart=on-failure
Environment=PYTHONHOME=/opt/wapt
Environment=PYTHONPATH=/opt/wapt
Environment=VIRTUAL_ENV=/opt/wapt
LimitNOFILE=32768

[Install]
WantedBy=multi-user.target
